# host/systemd/docker-cleanup.service
[Unit]
Description=Docker container and image cleanup service

[Service]
Type=oneshot
ExecStart=/opt/bin/docker-cleanup
