# host/systemd/docker-check-image-update.service
[Unit]
Description=Docker check image update service

[Service]
Type=oneshot
EnvironmentFile=<%= getenv!(:parasite_config_directory) %>/env/systemd.env
ExecStart=/opt/bin/docker-check-pull-all
