# unit/docker-check-image-update.timer
[Unit]
Description=Docker check image update timer

[Timer]
OnBootSec=5min
OnUnitInactiveSec=2min
Unit=docker-check-image-update.service
