# files/host/systemd/docker-check-image-update.service
[Unit]
Description=Docker check image update service

[Service]
Type=oneshot
ExecStart=/opt/bin/docker-check-pull-all
